# pi/kanp-collector.service
[Unit]
Description=KANP traffic collector (airplanes.live -> SQLite)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kanp
Group=kanp
ExecStart=/usr/bin/python3 /opt/kanp/pi/collector.py
Restart=always
RestartSec=15

# Tune here instead of editing collector.py:
#Environment=KANP_SOURCE=http://127.0.0.1/skyaware/data/aircraft.json
#Environment=KANP_POLL_SECONDS=15
#Environment=KANP_RADIUS_NM=60
#Environment=KANP_RETENTION_DAYS=365
#Environment=KANP_MAX_DB_MB=8000

# hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/kanp
PrivateTmp=true

[Install]
WantedBy=multi-user.target
