# pi/kanp-export.service
[Unit]
Description=KANP traffic snapshot exporter (SQLite -> traffic-data branch)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=kanp
Group=kanp
ExecStart=/usr/bin/python3 /opt/kanp/pi/exporter.py

#Environment=KANP_EXPORT_DIR=/var/lib/kanp/traffic-data
# Track simplification tolerance in nm: fixes are dropped only where they don't
# change a track's shape (no corner-cutting). Lower = more fidelity/bigger files,
# higher = smaller files. 0.03 nm (~180 ft) is visually lossless.
#Environment=KANP_SIMPLIFY_NM=0.03

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/kanp
PrivateTmp=true
